# deploy/systemd/rabbit-quant-deploy.timer
[Unit]
Description=Check Rabbit Quant production updates every minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
RandomizedDelaySec=10
Persistent=true
Unit=rabbit-quant-deploy.service

[Install]
WantedBy=timers.target
