# deploy/systemd/rabbit-quant-backup.timer
[Unit]
Description=Back up Rabbit Quant production data every day

[Timer]
OnCalendar=*-*-* 03:30:00 Asia/Shanghai
RandomizedDelaySec=15min
Persistent=true
Unit=rabbit-quant-backup.service

[Install]
WantedBy=timers.target
